# Infrastructure/certbot/certificate-manager.service
[Unit]
Description=Certificate Manager Daemon
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user

# Writable state directory under /var/lib
StateDirectory=certificate-manager
WorkingDirectory=%S

# Runtime directory for lock file
RuntimeDirectory=certificate-manager
RuntimeDirectoryMode=0755

# Set the log file path to match the certificate-manager directory
Environment=LOG_FILE=/var/log/certificate-manager/certificate-manager.log

ExecStart=/usr/bin/flock -n %t/certificate-manager/instance.lock /usr/local/bin/certificate-manager.sh --daemon

Restart=always
RestartSec=10

SyslogIdentifier=certificate-manager
StandardOutput=journal
StandardError=journal

LockPersonality=yes
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

LogsDirectory=certificate-manager

[Install]
WantedBy=multi-user.target
