# infrastructure/terraform/modules-single-az/compute/userdata/leader-manager/leader-manager.service
[Unit]
Description=Swarm Leader Lease Maintainer
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
Group=ec2-user
EnvironmentFile=/etc/leader-manager.env

# Runtime directory for lock file
RuntimeDirectory=leader-manager
RuntimeDirectoryMode=0755

# Log directory (systemd-managed, accessible even with ProtectSystem=strict)
LogsDirectory=leader-manager
LogsDirectoryMode=0755

ExecStart=/usr/local/bin/leader-manager.sh

# Script logs to stdout/stderr; systemd captures to logfile
StandardOutput=append:/var/log/leader-manager/leader-manager.log
StandardError=append:/var/log/leader-manager/leader-manager.log
SyslogIdentifier=swarm-leader

LockPersonality=yes
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# Explicitly allow writes to log directory (required with ProtectSystem=strict)
ReadWritePaths=/var/log/leader-manager
